# test/page/http-echo-headers@.service
[Unit]
Description=http-echo-headers

[Service]
ExecStart=-/home/alex/apps/http-echo-headers/http-echo-headers
StandardInput=socket
