# test/page/http-echo-headers.socket
[Unit]
Description=http-echo-headers socket

[Socket]
ListenStream=3008
Accept=yes

[Install]
WantedBy=sockets.target
